# systemd/regwatch-action.socket
[Unit]
Description=Socket for dispatching regwatch actions

[Socket]
ListenStream=127.0.0.1:13000
Accept=true

[Install]
WantedBy=sockets.target
