# systemd/regwatch.service
[Unit]
Description=Watch and warn about failed SIP registration
Wants=network.target
Requires=postfix.service
After=network.target postfix.service regwatch-action.socket

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/bin/regwatch.py

[Install]
WantedBy=multi-user.target
